# pxe/pxe-update-engine.service
[Service]
Type=oneshot
ExecStart=/usr/lib/coreos/pxe_update_engine
