# coreos-base/oem-rackspace/files/systemd/nova-agent.service
[Service]
ExecStart=/usr/share/oem/usr/bin/run-nova-agent.sh

[Install]
WantedBy=multi-user.target
