# coreos-base/oem-rackspace/files/systemd/setup-networking.service
[Service]
Type=oneshot
ExecStart=/usr/share/oem/usr/bin/net.sh

[Install]
WantedBy=multi-user.target
