# coreos-base/motd-http/files/motd-http.service
[Service]
ExecStart=/usr/sbin/motd-http

[Install]
WantedBy=multi-user.target
