# app-admin/systemd-rest/files/systemd-rest.service
[Service]
ExecStart=/usr/sbin/systemd-rest

[Install]
WantedBy=multi-user.target
