# app-admin/etcd/files/etcd.service
[Service]
ExecStart=/usr/sbin/etcd

[Install]
WantedBy=multi-user.target
