# net-wireless/bluez/files/bluetooth-init.d-r1
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-wireless/bluez/files/bluetooth-init.d-r1,v 1.1 2012/01/09 23:01:55 pacho Exp $

depend() {
	after coldplug
	need dbus localmount hostname
}

start() {
   	ebegin "Udev coldplug of bluetooth devices"
	udevadm trigger --subsystem-match=bluetooth --action=add
	eend $?
}

stop() {
	return 0
}
