# coreos-base/coreos-metadata/files/coreos-metadata-sshkeys@.service
[Unit]
Description=CoreOS Metadata Agent (SSH Keys)

[Service]
Type=oneshot
ExecStart=/usr/bin/coreos-metadata --cmdline --ssh-keys=%i

[Install]
DefaultInstance=core
RequiredBy=multi-user.target
