# app-crypt/go-tspi/files/tpmd.service
[Unit]
Description=TPM remote access daemon
Requires=tcsd.service

[Service]
User=tss
ExecStart=/bin/tpmd

[Install]
WantedBy=multi-user.target
