# app-emulation/xenstore/files/proc-xen.mount
[Unit]
Before=local-fs.target
DefaultDependencies=false
ConditionVirtualization=xen

[Mount]
What=xenfs
Where=/proc/xen
Type=xenfs
