# sys-fs/udev/files/50-compat_firmware.rules
# compat_firmware-class requests, copies files into the kernel
SUBSYSTEM=="compat_firmware", ACTION=="add", RUN+="compat_firmware.sh"
