# app-arch/torcx/files/torcx.target
[Unit]
Description=Verify torcx succeeded
DefaultDependencies=no
AssertPathExists=/run/metadata/torcx

[Install]
WantedBy=basic.target
