# app-admin/coreinit/files/coreinit.service
[Unit]
Description=coreinit

[Service]
ExecStart=/usr/bin/coreinit

[Install]
WantedBy=multi-user.target
