# net-misc/ntp/files/ntpdate.service
[Unit]
Description=Set time via NTP using ntpdate
After=network.target nss-lookup.target
Before=time-sync.target
Wants=time-sync.target

[Service]
Type=oneshot
Environment="SERVER=0.pool.ntp.org 1.pool.ntp.org 2.pool.ntp.org 3.pool.ntp.org"
ExecStart=/usr/sbin/ntpdate -b -u $SERVER
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
