# net-misc/modemmanager-next/files/modemmanager.conf
# Copyright (c) 2012 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# modemmanager
#
# starts modemmanager once udev has been started.

# Log level: ERR, WARN, INFO, DEBUG
env MM_LOGLEVEL=INFO

start on starting system-services
stop on stopping system-services

respawn

exec ModemManager --log-level="${MM_LOGLEVEL}"
