# app-admin/etcd-client/files/etcd-client.service
[Service]
ExecStart=/usr/sbin/etcd-client

[Install]
WantedBy=multi-user.target
