# app-emulation/docker/files/docker.service
[Service]
Type=forking
ExecStartPre=/bin/mount --make-rprivate /

# HACK: for some reason docker is crashing when exiting a container when being
# supervised by systemd. Fork out for now. Real solution should be:
# ExecStart=/usr/bin/docker -d -D

ExecStart=/bin/bash -c "/usr/bin/nohup /usr/bin/docker -d -D &"

[Install]
WantedBy=multi-user.target
