# app-emulation/docker/files/docker.service
[Service]
ExecStart=/usr/bin/docker -d

[Install]
WantedBy=multi-user.target
