# app-emulation/docker/files/early-docker.socket
[Unit]
Description=Early Docker Socket for the API

[Socket]
ListenStream=/var/run/early-docker.sock

[Install]
WantedBy=sockets.target
