# coreos-base/coreos-cloudinit/files/cloud-config.target
[Unit]
Requires=coreos-cloudinit-oem.service
After=coreos-cloudinit-oem.service

[Install]
WantedBy=multi-user.target
