# net-misc/modemmanager/files/modemmanager.conf
# Copyright (c) 2011 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# modemmanager
#
# starts modemmanager once udev has been started.

start on starting system-services

respawn

exec modem-manager
