# sys-apps/ignition/files/ignition-prepivot.target
[Unit]
Requires=ignition@prepivot.service
After=ignition@prepivot.service
