# dev-db/etcd/files/etcd.service
[Service]
ExecStartPre=-/bin/mkdir -p /var/run/etcd
ExecStartPre=/usr/bin/coreos-c10n
ExecStart=/usr/bin/etcd-bootstrap

[Install]
WantedBy=multi-user.target
