# dev-db/etcd/files/etcd.service
[Service]
ExecStart=/usr/bin/etcd-bootstrap

[Install]
WantedBy=multi-user.target
