# dev-db/etcd/files/etcd.service
[Unit]
Description=etcd
Wants=coreos-c10n.service
After=coreos-c10n.service

[Service]
User=etcd
PermissionsStartOnly=true
ExecStartPre=/usr/bin/etcd-pre-exec
ExecStart=/usr/bin/etcd-bootstrap

[Install]
WantedBy=multi-user.target
