# app-admin/fleet/files/fleet.service
[Unit]
Description=fleet
Wants=etcd.service
After=etcd.service

[Service]
ExecStart=/usr/bin/fleet
Restart=always
RestartSec=10s

[Install]
WantedBy=multi-user.target
