# coreos-base/oem-service/files/oem-service.service
[Service]
Type=oneshot
ExecStart=/usr/bin/run-oem.sh

[Install]
WantedBy=multi-user.target
